# ansible/templates/ai-translation.service.j2
[Unit]
Description=Simple AI Translation Service
After=network.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_group }}
WorkingDirectory={{ app_dir }}/ai-service
Environment=MODEL_NAME={{ model_name }}
Environment=MODEL_LOAD_MAX_RETRIES={{ model_load_max_retries }}
Environment=MODEL_LOAD_RETRY_DELAY_SEC={{ model_load_retry_delay_sec }}
Environment=HF_HOME={{ ai_cache_dir }}
Environment=TRANSFORMERS_CACHE={{ ai_cache_dir }}/transformers
ExecStart={{ ai_venv_path }}/bin/gunicorn --bind {{ ai_bind_host }}:{{ ai_port }} --workers {{ ai_workers }} --timeout 300 app:app
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ app_dir }}/ai-service {{ ai_cache_dir }}

[Install]
WantedBy=multi-user.target
